# .docker/release/files/etc/crontabs/root
# Run logrotate hourly
0 * * * * logrotate --state /tmp/logrotate.state /etc/logrotate.d/snappymail
